# templates/supervisorconf.ini
[program:php_{{ item.name }}]
command=/usr/bin/php-fpm -F -y /etc/php/php-fpm/{{ item.name }}.conf
user=php_{{ item.name }}
autorestart=true
redirect_stderr=true
